# supervisord.conf
[supervisord]
nodaemon=true

[program:nginx]
command=/usr/sbin/nginx -g 'daemon off;'
autostart=true
autorestart=true
stdout_logfile=/tmp/api_stdout.log  ; Redirect logs to files instead of stdout
stderr_logfile=/tmp/api_stderr.log
stdout_logfile_maxbytes=0           ; Disable log rotation
stderr_logfile_maxbytes=0

[program:api]
command=uvicorn hackathon.server.server:app --host 0.0.0.0 --port 3000
directory=/app
autostart=true
autorestart=true
stdout_logfile=/tmp/api_stdout.log  ; Redirect logs to files instead of stdout
stderr_logfile=/tmp/api_stderr.log
stdout_logfile_maxbytes=0           ; Disable log rotation
stderr_logfile_maxbytes=0
